# agent/deploy/pushlog-agent.service
[Unit]
Description=PushLog Agent — server log collector
Documentation=https://pushlog.ai/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pushlog-agent run
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pushlog-agent

# Hardening (ReadWritePaths=/var/run for Docker socket when using docker source)
# ProtectHome=false needed when file sources are in /home (e.g. ~/.pm2/logs)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/lib/pushlog-agent /var/run
ReadOnlyPaths=/etc/pushlog-agent
PrivateTmp=true

[Install]
WantedBy=multi-user.target
